# a/app-hello3.conf
start on starting app-hello3
stop on stopping app-hello3
respawn

env PORT=5302

setuid app

chdir /Users/cuihaoxin/go/src/github.com/hauson/goreman

exec ./loop.sh
