# a/app-hello1.conf
start on starting app-hello1
stop on stopping app-hello1
respawn

env PORT=5300

setuid app

chdir /Users/cuihaoxin/go/src/github.com/hauson/goreman

exec ./loop.sh
