# a/app-hello2.conf
start on starting app-hello2
stop on stopping app-hello2
respawn

env PORT=5301

setuid app

chdir /Users/cuihaoxin/go/src/github.com/hauson/goreman

exec ./loop.sh
